# battery-notifier.service
[Unit]
Description=Battery Low/High Percentage Notifier

[Service]
Restart=always
